# docker/s6-rc.d/nginx/run
#!/bin/sh -e

set -e

nginx -g "daemon off;"
