# docker/s6-rc.d/qcluster/run
#!/command/with-contenv sh

set -e

cd /var/www/redirect/backend

python3 manage.py qcluster
